# Telescope-Pi/etc/update-motd.d/1-telescope-pi
#!/bin/bash
echo "     _________"
echo "    / ======= \      Welcome to the"
echo "   / __________\      Telescope-Pi"
echo "  | ___________ |        shell"
echo "  | | -       | |"
echo "  | |         | |"
echo "  | |_________| |______________________"
echo "  \=____________/   marcocipriani01    )"
echo -e "  / \"\"\"\"\"\"\"\"\"\"\" \                     /"
echo " / ::::::::::::: \                =D-'"
echo "(_________________)"
echo ""
